# configs/com.user.tailscale-autoconnect.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.tailscale-autoconnect</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/a/.local/bin/tailscale-autoconnect.sh</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>StartInterval</key>
    <integer>600</integer>

    <key>LimitLoadToSessionType</key>
    <string>Aqua</string>

    <key>StandardOutPath</key>
    <string>/Users/a/.local/var/log/tailscale-autoconnect.launchd.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/a/.local/var/log/tailscale-autoconnect.launchd.log</string>
</dict>
</plist>
